# contrib/openwrt/ubavhid/files/ubavahid.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=80

start() {
	ubavahid br-lan
}

stop() {
	killall ubavahid
}
